# etc/systemd/log-vacuum.service
[Unit]
Description=Vacuum up logs

[Service]
Type=oneshot
ExecStart=/bin/bash journalctl --vacuum-size=20M
